# eval/fixtures-built/pipeline/patient-events/deploy/crontab
# Deployed schedule
0 4 * * * python3 jobs/daily_scoring.py
